# crontab
* * * * * /Users/nombauser/Desktop/GIT/my_git_repos/ETLWithBashScript/temperature_ETL.sh